# roles/system/templates/media-hdd.mount.j2
[Unit]
Description=Mount external HDD

[Mount]
What=/dev/sda1
Where={{ storage }}
Type=ext4
Options=noatime,nofail

[Install]
WantedBy=basic.target
